# EZ-WifiBroadcast/var/lib/dpkg/info/locate.list
/.
/etc
/etc/cron.daily
/etc/cron.daily/locate
/usr
/usr/bin
/usr/bin/updatedb.findutils
/usr/bin/locate.findutils
/usr/share
/usr/share/man
/usr/share/man/man5
/usr/share/man/man5/locatedb.5.gz
/usr/share/man/man1
/usr/share/man/man1/locate.findutils.1.gz
/usr/share/man/man1/updatedb.findutils.1.gz
/usr/share/doc
/usr/share/doc/locate
/usr/share/doc/locate/README.Debian
/usr/share/doc/locate/copyright
/usr/share/doc/locate/changelog.Debian.gz
/usr/share/doc/locate/changelog.gz
/usr/lib
/usr/lib/locate
/usr/lib/locate/bigram
/usr/lib/locate/code
/usr/lib/locate/frcode
/var
/var/cache
/var/cache/locate
