# EZ-WifiBroadcast/var/lib/dpkg/info/ser2net.conffiles
/etc/ser2net.conf
/etc/default/ser2net
/etc/init.d/ser2net
