# EZ-WifiBroadcast/var/lib/dpkg/info/ser2net.list
/.
/usr
/usr/bin
/usr/sbin
/usr/sbin/ser2net
/usr/share
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/ser2net.8.gz
/usr/share/doc
/usr/share/doc/ser2net
/usr/share/doc/ser2net/AUTHORS
/usr/share/doc/ser2net/README
/usr/share/doc/ser2net/copyright
/usr/share/doc/ser2net/examples
/usr/share/doc/ser2net/examples/ser2net.conf.gz
/usr/share/doc/ser2net/NEWS.gz
/usr/share/doc/ser2net/changelog.Debian.gz
/usr/share/doc/ser2net/changelog.gz
/etc
/etc/ser2net.conf
/etc/default
/etc/default/ser2net
/etc/init.d
/etc/init.d/ser2net
